# backup-daily.cron
#!/bin/bash

# Perform daily backups of system volumes

/etc/backup/backup-daily /home/ppzsb1 /data2/backup/home &> /data2/backup/logs/backup-daily-home.log
