# backup-monthly.cron
#!/bin/bash

# Perform monthly backup rotation

/etc/backup/backup-monthly ~/backup/projects
