# backup-weekly.cron
#!/bin/bash

# Perform weekly backup rotation

~/backup/bin/backup-weekly ~/backup/projects
