# backup-weekly.cron
#!/bin/bash

# Perform weekly backup rotation

/etc/backup/backup-weekly ~/backup/projects
