# backup-monthly.cron
#!/bin/bash

# Perform monthly backup rotation

~/backup/bin/backup-monthly ~/backup/projects
