# backup-monthly.cron
#!/bin/bash

# Perform monthly backup rotation for system volume backups

/etc/backup/backup-monthly /data2/backup/home
