# backup-weekly.cron
#!/bin/bash

# Perform weekly backup rotation for system volume backups

/etc/backup/backup-weekly /data2/backup/home
