# backup-daily.cron
#!/bin/bash

# Perform daily backups of system volumes

~ppzsb1/backup/backup-daily pppzsb2:projects ~/backup/projects &> ~/backup/backup-daily-projects.log
